# brains/cluster-brain.service
[Unit]
Description=ShaneBrain Brain Farm — Cluster Brain
After=network.target

[Service]
Type=simple
User=shanebrain
WorkingDirectory=/mnt/shanebrain-raid/shanebrain-core/mega-dashboard
ExecStart=/usr/bin/python3 brains/cluster-brain.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
